# door-manager.service
[Unit]
Description=Fox Lair Control System
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/door-manager
Environment="NODE_ENV=production"
Environment="PORT=3000"
ExecStart=/usr/bin/node server/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
